# supervisor.conf
[program:sqs-msg]
command=python3 /Users/divya/script/receive_sqs.py
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=/Users/divya/script/receive_sqs.log
[supervisord]